# free-pmx-no-shred/usr/lib/systemd/system/no-shred-flush.service
[Unit]
Description=free-pmx: NO SHRED flushing service
After=no-shred.service

[Service]
ExecStart=no-shred flush
